# contrib/whatsapp-daemon.service
[Unit]
Description=WhatsApp Baileys Daemon (singleton, manages auth_info and queue-processor)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/umbrel/mcp-servers/whatsapp-mcp-ts
Environment=NODE_ENV=production
Environment=LOG_LEVEL=info
ExecStart=/usr/bin/node --experimental-sqlite --experimental-strip-types --no-warnings /home/umbrel/mcp-servers/whatsapp-mcp-ts/src/daemon.ts
Restart=always
RestartSec=5
StandardOutput=append:%h/.local/share/whatsapp-daemon/stdout.log
StandardError=append:%h/.local/share/whatsapp-daemon/stderr.log

[Install]
WantedBy=default.target
